# com.amjad.talk.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.amjad.talk</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Use the Homebrew Python directly (not the venv symlink) so Python startup
             never reads pyvenv.cfg from ~/Documents/, which gets EDEADLK from launchctl.
             PYTHONPATH below makes venv packages visible. -->
        <string>/opt/homebrew/opt/python@3.12/bin/python3.12</string>
        <string>/Users/amjadhajireen/Documents/Claude Cowork/talk/start_talk.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/amjadhajireen/Documents/Claude Cowork/talk</string>

    <key>RunAtLoad</key>
    <true/>

    <!-- start_talk.py handles retry loop; launchctl restarts only if the launcher itself crashes -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Wait 60s between crash-restarts so macOS releases venv file locks -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONDONTWRITEBYTECODE</key>
        <string>1</string>
        <!-- Expose venv packages without triggering pyvenv.cfg reads at Python startup -->
        <key>PYTHONPATH</key>
        <string>/Users/amjadhajireen/Documents/Claude Cowork/talk/.venv/lib/python3.12/site-packages</string>
        <!-- Use macOS system SSL bundle instead of certifi (certifi is in ~/Documents, blocked by TCC) -->
        <key>SSL_CERT_FILE</key>
        <string>/etc/ssl/cert.pem</string>
        <key>REQUESTS_CA_BUNDLE</key>
        <string>/etc/ssl/cert.pem</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/amjadhajireen/Library/Logs/talk.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/amjadhajireen/Library/Logs/talk.log</string>
</dict>
</plist>
