# com.amjad.talk.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.amjad.talk</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/amjadhajireen/Documents/Claude Cowork/talk/.venv/bin/python</string>
        <string>/Users/amjadhajireen/Documents/Claude Cowork/talk/start_talk.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/amjadhajireen/Documents/Claude Cowork/talk</string>

    <key>RunAtLoad</key>
    <true/>

    <!-- start_talk.py handles retry loop; launchctl restarts only if the launcher itself crashes -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONDONTWRITEBYTECODE</key>
        <string>1</string>
        <!-- Use macOS system SSL bundle instead of certifi (certifi is in ~/Documents, blocked by TCC) -->
        <key>SSL_CERT_FILE</key>
        <string>/etc/ssl/cert.pem</string>
        <key>REQUESTS_CA_BUNDLE</key>
        <string>/etc/ssl/cert.pem</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/amjadhajireen/Library/Logs/talk.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/amjadhajireen/Library/Logs/talk.log</string>
</dict>
</plist>
