# com.amjad.talk.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.amjad.talk</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/amjadhajireen/Documents/Claude Cowork/talk/.venv/bin/python</string>
        <string>/Users/amjadhajireen/Documents/Claude Cowork/talk/talk.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/amjadhajireen/Documents/Claude Cowork/talk</string>

    <key>RunAtLoad</key>
    <true/>

    <!-- Restart on crash but not on clean exit (handles single-instance lock gracefully) -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Prevent rapid restart loops that cause pycache file-lock races -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Prevent pycache write contention on restart -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONDONTWRITEBYTECODE</key>
        <string>1</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/amjadhajireen/Library/Logs/talk.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/amjadhajireen/Library/Logs/talk.log</string>
</dict>
</plist>
